# crontab
*/5 * * * * /usr/bin/sh /root/remove-files.sh
