# crontab
*/5 * * * * /usr/bin/sh /root/gen-files.sh
